# src/mctech_crd/shell/bin/cosmic_service.service
[Unit]
Description=Listener for cosmic ray events.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/cosmic_service.sh

[Install]
WantedBy=multi-user.target
